# youtube-downloader.service
[Unit]
Description=YouTube Downloader - Liquid Glass UI
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/youtube-downloader
ExecStart=/usr/bin/python3 -m uvicorn app.main:app --host 0.0.0.0 --port 8081
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
